# autostart/.config/autostart/remap-keys.desktop
[Desktop Entry]
Type=Application
Exec=/home/calum/.xprofile
Hidden=false
X-GNOME-Autostart-enabled=true
Name=Apply Xmodmap
